# tdarr-server.service
[Unit]
Description=Tdarr Server - Transcoding application manager
Documentation=https://docs.tdarr.io/
After=network.target

[Service]
Type=simple
User=tdarr
Group=tdarr
WorkingDirectory=/opt/tdarr/server
ExecStart=/opt/tdarr/server/Tdarr_Server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tdarr /etc/tdarr

# Resource limits (adjust as needed)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

